# nginx/package/DEBIAN/postrm
#!/bin/bash
# set -e

systemctl stop nginx
systemctl disable nginx
rm -f /etc/nginx/sites-enabled/tooloop
rm -f /assets/logs/access.log
rm -f /assets/logs/error.log

# exit 0